# nixos-new/old/sway/.config/systemd/user/lxqt-policykit-agent.service
[Unit]
Description=LXQt PolicyKit Agent

[Service]
ExecStart=/usr/bin/lxqt-policykit-agent

[Install]
WantedBy=sway-session.target

